# scheduler/crontab
0 1 1 * 1 sh /home/bates/repositorio/big_data/one_piece/run.sh